# auto/com.traderloadings.nightshift.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.traderloadings.nightshift</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/caffeinate</string>
    <string>-i</string>
    <string>/Users/gazz/Desktop/TraderLoadingsLOCALE/auto/run.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>NIGHTSHIFT_SCHEDULED</key><string>1</string>
    <key>HOME</key><string>/Users/gazz</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict><key>Hour</key><integer>1</integer><key>Minute</key><integer>0</integer></dict>
  <key>StandardOutPath</key><string>/Users/gazz/Desktop/TraderLoadingsLOCALE/auto/logs/launchd.out.log</string>
  <key>StandardErrorPath</key><string>/Users/gazz/Desktop/TraderLoadingsLOCALE/auto/logs/launchd.err.log</string>
  <key>RunAtLoad</key><false/>
</dict>
</plist>
